# docker-mynet.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/docker-mynet.sh

[Install]
WantedBy=default.target
